# scripts/deployment/aqea-master.service
[Unit]
Description=AQEA Distributed Extractor - Master Coordinator
After=network.target

[Service]
Type=simple
User=aqea
Group=aqea
WorkingDirectory=/opt/aqea-distributed-extractor
ExecStart=/opt/aqea-distributed-extractor/aqea-venv/bin/python -m src.main start-master --language de --workers 4 --source wiktionary --port 8080
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aqea-master
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target 