# scripts/deployment/aqea-worker@.service
[Unit]
Description=AQEA Distributed Extractor - Worker %i
After=network.target aqea-master.service
Requires=aqea-master.service

[Service]
Type=simple
User=aqea
Group=aqea
WorkingDirectory=/opt/aqea-distributed-extractor
ExecStart=/opt/aqea-distributed-extractor/aqea-venv/bin/python -m src.main start-worker --worker-id worker-%i --master-host localhost --master-port 8080
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aqea-worker-%i
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target 